# runs/systemd/tfmx-tei-machine.service.in
[Unit]
Description=tfmx TEI machine service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@
Environment=HOME=%h
Environment=PATH=@PATH@
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env bash -lc 'GPUS="${TEI_SERVICE_GPUS:-$(nvidia-smi --query-gpu=index --format=csv,noheader,nounits 2>/dev/null | paste -sd, -)}"; exec "@PYTHON@" -m tfmx.teis.cli machine run --auto-start --compose-gpus "$GPUS" --on-conflict replace --perf-track'
ExecStopPost=/usr/bin/env bash -lc 'GPUS="${TEI_SERVICE_GPUS:-$(nvidia-smi --query-gpu=index --format=csv,noheader,nounits 2>/dev/null | paste -sd, -)}"; cd "@REPO_ROOT@" && "@PYTHON@" -m tfmx.teis.cli compose down -g "$GPUS" || true'
Restart=always
RestartSec=10
TimeoutStartSec=0
TimeoutStopSec=180

[Install]
WantedBy=default.target